# roles/application/cups-print-server/files/resources/services/cupsrestartservice/cups.enable.printers.timer
[Unit]
Description=Run enable printers frequently to ensure connection difficulties are remedied.

[Timer]
OnBootSec=15min
OnUnitActiveSec=300

[Install]
WantedBy = timers.target